# telerob_deployment/ansible/roles/setup_dds_prerequisites/files/50-loopback-if.rules
# Enable multicast for the loopback interface
ACTION=="add", KERNEL=="lo", SUBSYSTEM=="net", RUN+="/sbin/ip link set dev %k multicast on"
