# telerob_deployment/ansible/roles/system_services/files/udev/99-uart-imu-config.rules
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyTHS1", RUN+="stty -F /dev/%k -crtscts"
